# supervisor/staging.conf
[program:queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/trip2/current/artisan queue:work
autostart=true
autorestart=true
user=tripikas
numprocs=1
redirect_stderr=true
stdout_logfile=syslog